# server.conf
[program:myserver]
command=gunicorn run:app -w 4 -b 0.0.0.0:80
directory=/home/ubuntu/src/happening/

[supervisord]
logfile=/home/ubuntu/src/happening/supervisord_server.log
loglevel=debug
user=root
